# backup/protonvpn-autostart/.config/autostart/com.protonvpn.www.desktop
[Desktop Entry]
Type=Application
Name=Proton VPN
Exec=/usr/bin/flatpak run com.protonvpn.www --start-minimized
